# pok-crawler.service
[Unit]
Description=Pok Webcrawler

[Service]
Type=simple
ExecStart=/var/pok/pok-crawler/start-crawler.sh

[Install]
WantedBy=multi-user.target